# etc/init/redis-dns.conf
# redis-dns - run redis-dns as service in ubuntu

description "redis-dnsas a service"
author "Jonas Colmsjö <jonas@gizur.com>"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Essentially lets upstart know the process will detach itself to the background
#expect fork

# Run before process
pre-start script
#    [ -d /var/run/docker-openerp ] || mkdir -p /var/run/docker-openerp
    echo "Starting redis-dns..."
end script

# Start the process
exec /usr/bin/node /usr/lib/node_modules/redis-dns/server.js > /var/log/redis-dns.log 2>&1
